# backend/supervisor-workers.conf
[program:laravel-claude-events]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/live-agents.tams.codes/backend/artisan queue:work redis-high-throughput --queue=claude-events --sleep=1 --tries=3 --max-time=3600 --memory=128
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/live-agents.tams.codes/backend/storage/logs/claude-events-worker.log
stopwaitsecs=3600

[program:laravel-agent-execution]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/live-agents.tams.codes/backend/artisan queue:work database --queue=agent-execution --sleep=3 --tries=3 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/live-agents.tams.codes/backend/storage/logs/agent-execution-worker.log
stopwaitsecs=3600

[program:laravel-messages]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/live-agents.tams.codes/backend/artisan queue:work database --queue=messages --sleep=2 --tries=3 --max-time=3600 --memory=64
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/live-agents.tams.codes/backend/storage/logs/messages-worker.log
stopwaitsecs=3600

[program:laravel-commands]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/live-agents.tams.codes/backend/artisan queue:work database --queue=commands --sleep=3 --tries=3 --max-time=3600 --memory=128
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/live-agents.tams.codes/backend/storage/logs/commands-worker.log
stopwaitsecs=3600

[program:laravel-webhooks]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/live-agents.tams.codes/backend/artisan queue:work database --queue=webhooks --sleep=5 --tries=5 --max-time=3600 --memory=64
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/live-agents.tams.codes/backend/storage/logs/webhooks-worker.log
stopwaitsecs=3600

[program:laravel-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/live-agents.tams.codes/backend/artisan queue:work database --queue=default --sleep=3 --tries=3 --max-time=3600 --memory=128
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/live-agents.tams.codes/backend/storage/logs/default-worker.log
stopwaitsecs=3600

[group:laravel-workers]
programs=laravel-claude-events,laravel-agent-execution,laravel-messages,laravel-commands,laravel-webhooks,laravel-default
priority=999
