# backend/docker/supervisor/openclaw-watcher.conf
[program:openclaw-watcher]
process_name=%(program_name)s
command=php /var/www/html/artisan openclaw:watch --interval=2
autostart=true
autorestart=true
startretries=3
startsecs=5
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/openclaw-watcher.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=10
stopsignal=SIGTERM
